# mini_os.service
[Unit]
Description=Mini OS Pico-8 Launcher
After=network.target

[Service]
Type=simple
# Adjust the paths below if the project lives elsewhere.
# This example assumes it resides under the service user's home
# directory at ~/pico.
ExecStart=/usr/bin/python3 %h/pico/main.py
WorkingDirectory=%h/pico
Restart=always
User=pi
Group=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
